# etc/periodic/hourly/cleanup
#!/bin/sh

ruby /usr/local/bin/cleanup.rb -l 95G /srv/backstore
