# etc/services.d/caddy/run
#!/usr/bin/execlineb -P

caddy run --config /etc/caddy/Caddyfile --adapter caddyfile

