# etc/services.d/ssh/finish
#!/usr/bin/execlineb -S0

# Ensure the container stops if sshd fails for same reason.
s6-svscanctl -t /var/run/s6/services
