# etc/services.d/ssh/run
#!/usr/bin/execlineb -P

/usr/sbin/sshd -Dep 1984

